# systemd/update-check.service
[Unit]
Description=Check for updates and report results

[Service]
Type=oneshot
ExecStart=update-check.sh