# systemd/update-check.timer
[Unit]
Description=Timer for update checker

[Timer]
OnBootSec=10min
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target